# update.cron
10 */2 * * * /usr/local/bin/wayback-rpki update
